# templates/default/kafka.upstart.conf.erb
description "Kafka"

limit nofile 32768 32768

start on <% if @local_zoo %>started zookeeper<% else %>[2345]<% end %>
stop on <% if @local_zoo %>stopped zookeeper<% else %>[!12345]<% end %>

respawn
respawn limit 2 5

umask 007

kill timeout 300

setuid <%= @user %>
setgid <%= @group %>

exec <%= @install_directory %>/kafka-<%= @version %>/bin/kafka-server-start.sh <%= @install_directory %>/kafka-<%= @version %>/config/server-<%= @n %>.properties;