# templates/default/kafka.upstart.conf.erb
description "Kafka"

limit nofile 32768 32768

start on [2345]
stop on [!2345]

respawn
respawn limit 2 5

umask 007

kill timeout 300

setuid <%= @user %>
setgid <%= @group %>

exec <%= @install_directory %>/kafka_<%= @scala_version %>-<%= @version %>/bin/kafka-server-start.sh <%= @install_directory %>/kafka_<%= @scala_version %>-<%= @version %>/config/server-<%= @n %>.properties;

post-start script
  PID=`status kafka-<%= @n %> | egrep -oi '([0-9]+)$' | head -n1`
  echo $PID > <%= @install_directory %>/kafka-<%= @n %>.pid
end script

post-stop script
  rm -f <%= @install_directory %>/kafka-<%= @n %>.pid
end script
